# deploy/uberspace/sensor_board.ini
; Uberspace supervisord service for Sensor Board.
; Install: copy to ~/etc/services.d/sensor_board.ini on the Uberspace host,
; edit the environment line (secrets + domain), then:
;   supervisorctl reread && supervisorctl update
;
; Paths use %(ENV_HOME)s so this file is account-agnostic.
[program:sensor_board]
directory=%(ENV_HOME)s/sensor_board
command=%(ENV_HOME)s/sensor_board/.venv/bin/gunicorn -c conf.py

; All configuration (API_KEYS, ROOT_PATH, BASE_URL, …) lives in
; ~/sensor_board/.env — the app loads it itself. Do NOT put comma-containing
; values like API_KEYS here: supervisord's environment= splits on commas.

autostart=yes
autorestart=yes
startsecs=30
stopwaitsecs=10
stopsignal=TERM
