# deploy/uberspace/sensor_board.ini
; Uberspace supervisord service for Sensor Board.
; Install: copy to ~/etc/services.d/sensor_board.ini on the Uberspace host,
; edit the environment line (secrets + domain), then:
;   supervisorctl reread && supervisorctl update
;
; Paths use %(ENV_HOME)s so this file is account-agnostic.
[program:sensor_board]
directory=%(ENV_HOME)s/sensor_board
command=%(ENV_HOME)s/sensor_board/.venv/bin/gunicorn -c conf.py

; --- Edit these before starting ---
; API_KEYS:  comma-separated; any one is accepted as the x-api-key header.
; ROOT_PATH: URL prefix the UI is mounted under ("" = domain root).
; BASE_URL:  public https URL of the dashboard UI (include ROOT_PATH).
environment=API_KEYS="CHANGE_ME_KEY_1,CHANGE_ME_KEY_2",ROOT_PATH="/dashboard",BASE_URL="https://CHANGE_ME_DOMAIN/dashboard",APP_TITLE="Sensor Board",BRAND="Sensor Board"

autostart=yes
autorestart=yes
startsecs=30
stopwaitsecs=10
stopsignal=TERM
